# gps-ntp-server.service
[Unit]
Description=GPS NTP Time Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gps-ntp-server
ExecStart=/usr/bin/python3 /opt/gps-ntp-server/gps_ntp_server.py --gps-port /dev/ttyUSB0 --gps-baud 9600 --web-port 5000 --ntp-port 123
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening (comment out if issues arise)
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gps-ntp-server

[Install]
WantedBy=multi-user.target
