# gps-ntp-server.service
[Unit]
Description=GPS NTP Time Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gps-ntp-server
# Use the virtual environment Python to ensure dependencies are available
ExecStart=/opt/gps-ntp-server/venv/bin/python3 /opt/gps-ntp-server/gps_ntp_server.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening (optional - uncomment if needed, but may cause issues with device access)
#PrivateTmp=yes
#NoNewPrivileges=true
#ProtectSystem=strict
#ProtectHome=true
#ReadWritePaths=/opt/gps-ntp-server

[Install]
WantedBy=multi-user.target
